# src/scripts/backup.cron
0 5 * * * /bin/sh /scripts/backup.sh >> /var/log/cron.log 2>&1